# raspi-botutils.service
[Unit]
Description=Raspi Monitor Bot - Ultimate Edition
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/raspi-botutils
ExecStart=/usr/bin/python3 /opt/raspi-botutils/run_bot.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/raspi-botutils/.env
StandardOutput=journal
StandardError=journal
WatchdogSec=60s

[Install]
WantedBy=multi-user.target
